# tools/openrc/sbdaemon
#!/sbin/runscript
# Short-Description: Shadowbound deamon
# Description:       Shadowbound daemon used to start the server and keep it updated

source /etc/shadowbound/shadowbound.cfg

NAME="StarboundServer"
LOGFILE="${logdir}/${NAME}.log"
DAEMON="/usr/bin/shadowbound"

depend(){
    need net
}

start(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Starting Shadowbound daemon"
    ulimit -n 100000
    if [ "$INSTANCE" != "$RC_SVCNAME" ];  then
        "$DAEMON" start "@${INSTANCE}"
    else
        "$DAEMON" start "@all"
    fi
    eend $?
}

stop(){
    INSTANCE="${RC_SVCNAME#*.}"
    ebegin "Stopping Shadowbound daemon"
    if [ "$INSTANCE" != "$RC_SVCNAME" ]; then
        "$DAEMON" stop "@${INSTANCE}"
    else
        "$DAEMON" stop "@all"
    fi
    eend $?
}