# tools/systemd/shadowbound@.service
[Unit]
Description=Daemon to start an starbound server instance

[Service]
User=steam
ExecStart=/usr/bin/shadowbound run
Type=simple
