# tools/systemd/shadowbound@.service
[Unit]
Description=Daemon to start an starbound server instance

[Service]
ExecStart=/usr/bin/shadowbound useconfig %i run
Type=simple
